# source/build/package-templates/hdn-nginx/hdn-nginx/etc/init/hdn-nginx.conf
# © Crown Copyright 2013
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

description "hdn-nginx"
author "Raphael Cohn <raphael.cohn@stormmq.com>"

start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!2345]
 
expect fork

# 2 min timeout - possibly not enough
kill timeout 120

# Graceful shutdown rather than quick
# Not supported on lucid
kill signal QUIT

limit core unlimited unlimited
limit cpu unlimited unlimited
limit data unlimited unlimited
limit fsize unlimited unlimited
limit memlock unlimited unlimited
limit msgqueue unlimited unlimited
limit nice unlimited unlimited
limit nofile 1048576 1048576
limit nproc unlimited unlimited
limit rss unlimited unlimited
limit rtprio unlimited unlimited
limit sigpending unlimited unlimited
limit stack unlimited unlimited

nice -4
oom never

respawn
respawn limit 5 60
 
pre-start script
	#!/bin/sh -e
	/usr/sbin/hdn-nginx -t
	if [ $? -ne 0 ]; then
		exit $?
	fi
end script
 
exec /usr/sbin/hdn-nginx
