# source/build/package-templates/hdn-crds-registry-server/hdn-crds-registry-server/etc/init/hdn-crds-registry-server.conf
# © Crown Copyright 2013
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

description "hdn-crds-registry-server"
author "Raphael Cohn <raphael.cohn@stormmq.com>"

# needs review?
start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]

console none

# 2 min timeout - possibly not enough
kill timeout 120

limit core unlimited unlimited
limit cpu unlimited unlimited
limit data unlimited unlimited
limit fsize unlimited unlimited
limit memlock unlimited unlimited
limit msgqueue unlimited unlimited
limit nice unlimited unlimited
# unlimited causes failure of the job to start as it exceeds Linux's compiled in maximum
limit nofile 1048576 1048576
limit nproc unlimited unlimited
limit rss unlimited unlimited
limit rtprio unlimited unlimited
limit sigpending unlimited unlimited
limit stack unlimited unlimited

nice -4
oom never

respawn
respawn limit 5 60

# script is run as #!/bin/sh -e (which has different syntax to bash)
script
	#!/bin/sh -e
	defaultFile=/etc/default/hdn-crds-registry-server
	if [ -f "$defaultFile" ]; then
		profile="$(head -q -n1 "$defaultFile")"
		exec /usr/sbin/hdn-crds-registry-server $profile
	else
		exec /usr/sbin/hdn-crds-registry-server
	fi
end script
