# source/build/package-templates/hdn-jstatd/hdn-jstatd/etc/init/hdn-jstatd.conf
# © Crown Copyright 2013
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
# Unless required by applicable law or agreed to in writing, software
# distributed under the License is distributed on an "AS IS" BASIS,
# WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
# See the License for the specific language governing permissions and
# limitations under the License.

description "hdn-jstatd"
author "Raphael Cohn <raphael.cohn@stormmq.com>"

# needs review?
start on (local-filesystems and net-device-up IFACE=eth0)
stop on runlevel [!12345]

console none

kill timeout 5

nice 15

respawn
respawn limit 5 60

# script is run as #!/bin/sh -e (which has different syntax)
script
	#!/bin/sh -e
	exec /usr/sbin/hdn-jstatd
end script
